# assignment_04/999-keyboard-add.rules
ACTION=="add", ENV{ID_BUS}=="usb", ENV{ID_INPUT_KEYBOARD}=="?*", RUN+="/sbin/insmod /root/little-penguin/assignment_04/main.ko"
